# ssl-checker/systemd/ssl-checker.timer
[Unit]
Description=Run ssl-checker daily
Requires=ssl-checker.service

[Timer]
OnCalendar=daily
AccuracySec=1h

[Install]
WantedBy=timers.target